# src/brain/_assets/scripts/brain-brief-mac.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- brain-daily-brief — macOS launchd LaunchAgent (UX-02)
     Install: BRAIN_VAULT=/path/to/vault bash scripts/install-brief-mac.sh
     Uninstall: launchctl unload ~/Library/LaunchAgents/com.brainiac.nightly.<id>.plist
                rm ~/Library/LaunchAgents/com.brainiac.nightly.<id>.plist
     Logs: ~/Library/Logs/brain/brief-YYYY-MM-DD.log
     Task name: com.brainiac.nightly.<id> (per-vault; <id> = 8-hex vault slug, rendered by install-brief-mac.sh)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>LABEL_PLACEHOLDER</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>SCRIPTS_DIR/brain-brief.sh</string>
    </array>

    <!-- BRAIN_VAULT and BRAIN_AUDIT_KEY_PEM are substituted by install-brief-mac.sh.
         The signing key is read from the macOS Keychain at install time and written
         here; it is never stored in the source tree. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAIN_VAULT</key>
        <string>VAULT_PATH</string>
        <key>BRAIN_AUDIT_KEY_PEM</key>
        <string>AUDIT_KEY_PEM_PLACEHOLDER</string>
        <key>BRAIN_LOG_DIR</key>
        <string>HOME_DIR/.brain/logs</string>
        <!-- Workspace sweep sources (WSP-01, optional): os.pathsep-separated
             folders whose SETTLED top-level files auto-flow into
             <vault>/inbox/. `path=N` sets a per-dir age gate in days;
             `path=0` = capture-inbox mode (same-day, 15-min write guard).
             Empty = sweep disabled. Set per vault at install time. -->
        <key>BRAIN_WORKSPACE_SWEEP_DIRS</key>
        <string>SWEEP_DIRS_PLACEHOLDER</string>
    </dict>

    <!-- HOURLY (owner decision 2026-07-11: ingestion is frequent — a
         captured document is searchable within the hour). Every firing runs
         the incremental/idempotent work (sweep + ingest + drain + sync +
         snapshot + self-organization folds); weekly/monthly branches stay
         date-gated inside `brain maintain` via .brain/maintain-state.json,
         so higher frequency never multiplies the heavy branches. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>HOME_DIR/.brain/logs/launchd-brief-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_DIR/.brain/logs/launchd-brief-stderr.log</string>

    <!-- Do NOT run immediately on load — wait for the first scheduled time. -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
